# desktop/config/systemd/user/tmux.service.d/override.conf
[Service]
ExecStart=
ExecStart=/bin/sh -lc "/usr/bin/tmux new-session -d"

ExecStop=
ExecStop=/usr/bin/tmux kill-server
