# desktop/config/systemd/user/tmux.service.d/override.conf
[Service]
ExecStop=
ExecStop=/usr/bin/tmux kill-server
